# deploy/supervisor.conf
# ─── PoliticOS — Supervisor config ────────────────────────────────────────────
# Place in /etc/supervisor/conf.d/politicos.conf
# sudo supervisorctl reread && sudo supervisorctl update

# ── Laravel queue workers ─────────────────────────────────────────────────────
[program:politicos-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/politicos/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/politicos-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=60

# ── Next.js app server ────────────────────────────────────────────────────────
[program:politicos-nextjs]
process_name=%(program_name)s
command=node /var/www/politicos/resources/js/.next/standalone/server.js
directory=/var/www/politicos/resources/js
environment=NODE_ENV="production",PORT="3000",HOSTNAME="127.0.0.1"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/politicos-nextjs.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

# ── Group: restart both together ──────────────────────────────────────────────
[group:politicos]
programs=politicos-queue,politicos-nextjs
