# deploy/supervisor.conf
# ─── PoliticOS — Supervisor config ────────────────────────────────────────────
# Place in /etc/supervisor/conf.d/politicos.conf
# sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start politicos:*

# ── Workers generales: geo, análisis de mensajes, alertas, etc. ───────────────
# 3 procesos paralelos para absorber picos de chat
[program:politicos-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/politicos/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/politicos-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=60

# ── Worker de larga duración: merge de chunks de video (hasta 2 horas) ────────
# 1 proceso dedicado — timeout extendido, no compite con jobs cortos
[program:politicos-queue-video]
process_name=%(program_name)s
command=php /var/www/politicos/artisan queue:work redis --queue=video --sleep=5 --tries=1 --timeout=7200 --memory=512
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/politicos-queue-video.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=7200

# ── Next.js app server ────────────────────────────────────────────────────────
[program:politicos-nextjs]
process_name=%(program_name)s
command=node /var/www/politicos/resources/js/.next/standalone/server.js
directory=/var/www/politicos/resources/js
environment=NODE_ENV="production",PORT="3000",HOSTNAME="127.0.0.1"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/politicos-nextjs.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=10

# ── Group: reiniciar todo junto con "supervisorctl restart politicos:*" ────────
[group:politicos]
programs=politicos-queue,politicos-queue-video,politicos-nextjs
