# deploy/render/supervisord.conf
; ─────────────────────────────────────────────────────────────────────────────
; PoliticOS — supervisord para el contenedor de Render
; Corre php-fpm y nginx en primer plano; si uno muere, supervisord lo reinicia.
; NO hay workers de cola: la demo usa QUEUE_CONNECTION=sync (sin Redis).
; ─────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
